# roles/nextcloud_server/templates/nextcloud.cron.j2
#
# Generated by Ansible
#

# https://doc.owncloud.org/server/8.0/admin_manual/configuration_server/background_jobs_configuration.html
# 2015-08-23 LL


HOME=/tmp
2-59/5  *  *  *  *  {{ nextcloud.user }} /usr/bin/php -f {{ nextcloud.installdir }}/nextcloud/cron.php > /dev/null 2>&1

# 2,7,12 etc...
# https://stackoverflow.com/questions/12786410/run-cron-job-every-n-minutes-plus-offset